# etc/init/black.conf
start on startup

task

script
  if [ -f /forcefsck ]; then 
     exit 0
  fi 

  setterm -foreground black -store >/dev/tty1
  setterm -foreground black -store >/dev/tty7
  setterm -foreground black -store >/dev/tty9
  setterm -cursor off -store >/dev/tty1
  setterm -cursor off -store >/dev/tty7
  setterm -cursor off -store >/dev/tty9
end script

post-start script
  if [ -f /forcefsck ]; then
     exit 0
  fi
sleep 30
end script

