# etc/init/onshutdown.conf
description     "do this and that on shutdown or reboot"

start on runlevel [06]

task

script
# disable some startup skript left for compatibility
  /usr/sbin/update-rc.d -f vdr remove >/dev/null
  /usr/sbin/update-rc.d -f inputlirc remove >/dev/null

end script

