# etc/init/s3.conf
# Sends PC to sleep

start on init-s3 
stop on runlevel [!2345]

script
/etc/acpi/sleep.sh &
end script
