# etc/init/s4.conf
# Sends PC to hibernate

start on (init-s4 and vdr-stopped)
stop on runlevel [016]

script
/usr/sbin/pm-hibernate
end script
