# etc/init/svdrpd.conf
start on started vdr
stop on stopping vdr

script
. /usr/lib/vdr/config-loader

exec /usr/sbin/svdrpd localhost 6420 localhost $SVDRP_PORT
end script

