# etc/init/vdr-exit-other.conf
start on ( stopped vdr RESULT=failed EXIT_SIGNAL=?* or \
           stopped vdr RESULT=failed EXIT_STATUS!=[012] )

task

script
logger -t "vdr-crash" "vdr exit with exit status $EXIT_STATUS . Restarting"
end script


