# etc/init/vdr-exit-2.conf
start on stopped vdr RESULT=failed EXIT_STATUS=2

task

script
DATESTAMP=`date +%Y%m%d%M%S`
cp /tmp/vdr.log /var/log/vdr/vdr.log.$DATESTAMP
# check to see which plugin causes problems
end script

