# etc/init/white.conf

start on (stopped black or starting mountall-shell) 

task 

script 
setterm -clear > /dev/tty1
setterm -clear > /dev/tty7
setterm -clear > /dev/tty9
end script 

