# etc/udev/rules.d/40-sound-card-added.rules
SUBSYSTEM=="sound" , ACTION=="change", ENV{SOUND_INITIALIZED}="1", RUN+="/sbin/initctl emit --no-wait sound-card-added CARD=%p"
